# services/php/run
#!/bin/sh
exec 2>&1
source /etc/envvars
install -d -o librenms -g librenms -m 755 /run/php

exec /usr/sbin/php-fpm7 --nodaemonize
